# common-linux/bin/.config/systemd/user/battery-watch.service
[Unit]
Description=Battery threshold watcher

[Service]
Type=oneshot
Environment=LOW=20
Environment=CRIT=10
ExecStart=%h/.local/bin/battery-watch.sh
