# common-linux/bin/.config/systemd/user/home-osteiner-Nextcloud.automount
[Unit]
Description=Automount trigger for %h/Nextcloud

[Automount]
Where=%h/Nextcloud

[Install]
WantedBy=default.target
