# common-linux/hypr-asahi/.config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland session is up
